# archlive/airootfs/etc/systemd/system/install-ml.service
[Unit]
Description=Install ML/AI frameworks and NVIDIA stack
After=network-online.target
Wants=network-online.target
ConditionPathExists=!/var/lib/install-ml.done

[Service]
Type=oneshot
ExecStart=/usr/local/bin/install-ml.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
